# crontab
* * * * * /usr/bin/systemctl restart refreshscript.service